# packaging/init/photond.conf
start on (filesystem and net-device-up)

#wherever you want to store data
chdir /var/lib/photond
console log

script
    /usr/local/bin/photond -s map
end script
