# packaging/init/photond.conf
start on (filesystem and net-device-up)

#wherever you want to store data
chdir /var/lib/ortd
console log

script
    /usr/local/bin/ortd -s map
end script
