# packaging/root/usr/share/ort/init/ortd.conf
start on (filesystem and net-device-up and remote-filesystem)

chdir /var/lib/ort
console log

kill signal SIGINT
kill timeout 20

script 
    /usr/local/bin/ortd
end script
